# .config/systemd/user/zathura.service
[Unit]
 Description=Zathura
 Documentation=zathura

 [Service]
 ExecStart=/usr/bin/zathura
