# .config/autostart/Logitech Solaar.desktop
[Desktop Entry]
Exec=/usr/bin/solaar -w hide
Name=Logitech Solaar
Type=Application
Version=1.0
